# address-validator.service
[Unit]
Description=Address Validator (FastAPI)
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/address-validator
EnvironmentFile=/etc/address-validator/.env
EnvironmentFile=-/home/exedev/address-validator/.env
# Validation provider config (optional — defaults to none/unavailable):
#   VALIDATION_PROVIDER=usps
#   USPS_CONSUMER_KEY=<key from developer.usps.com>
#   USPS_CONSUMER_SECRET=<secret from developer.usps.com>
StateDirectory=address-validator
Environment=PYTHONPATH=/home/exedev/address-validator/src
ExecStart=/home/exedev/address-validator/.venv/bin/uvicorn address_validator.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
